# deploy/systemd/ask-insects.service
[Unit]
Description=Ask Insects hosted source service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=josh
WorkingDirectory=/home/josh/ask-insects
# ASK_INSECTS_TOKEN is loaded from this file.
EnvironmentFile=/home/josh/ask-insects/.env
ExecStart=/usr/bin/python3 -m askinsects.server --host 127.0.0.1 --port 8080 --artifact-dir /home/josh/ask-insects/artifacts/mosquito-v1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
